# packages/files/ets/inite.d/phoenix
#!/bin/sh /etc/rc.common

START=90

start() {
    logger "Phoenix Core Started"
}

stop() {
    logger "Phoenix Core Stopped"
}
